# etc/inittab
::sysinit:/etc/rc
::respawn:-/bin/sh
slog:unknown:/sbin/syslogd -n
